# {{cookiecutter.repostory_name}}/backups/cron.d/backup
0 3 * * * root /root/backup-db.sh >> /var/log/cron.log 2>&1
0 5 * * * root /root/rotate-local-backups.sh >> /var/log/cron.log 2>&1

