# {{cookiecutter.repostory_name}}/backups/cron.d/backup
0 3 * * * root /root/backup-db.sh > /proc/1/fd/1 2>&1
0 5 * * * root /root/rotate-local-backups.sh > /proc/1/fd/1 2>&1
